# start_scope.service
[Unit]
Description=Start Weave Scope on boot

[Service]
ExecStart=/home/administrador/start_scope.sh

[Install]
WantedBy=multi-user.target
